# cron/indexation.cron
15 3  * * 1-5 sudo php /var/www/mc3/bin/console indexation:start
